# service/sbds-streamer/run
#!/bin/bash

exec setuidgid nobody bash /app/sbds/storages/db/scripts/populate.sh
