# service/sbds-streamer/run
#!/bin/bash

exec setuidgid nobody bash populate
