# service/sbds-streamer/run
#!/bin/bash

exec setuidgid nobody bash /app/sbds/scripts/populate-sbds.sh
