# service/sbds-web/run
#!/bin/bash

exec setuidgid nobody dev-server
